# .config/systemd/user/conky-autostart.service
[Unit]
Description=Conky Startup (Wayland-safe with delay)
After=graphical-session.target

[Service]
Type=simple
ExecStart=/home/csantos/.config/conky/Maia/start.sh
